# apps/dootask/1.8.64/scripts/dootask-supervisor-cron.conf
[program:dootask-cron]
directory=/var/www/docker/crontab
command=/usr/sbin/cron -f
numprocs=1
autostart=true
autorestart=false
startretries=1
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
